# scheduler.cron
* * * * * curl -s 'http://localhost/api/v1/xapi-extract' >> /var/log/cron.log 2>&1
* * * * * /usr/local/bin/php /var/www/html/artisan schedule:run >> /dev/null 2>&1
